# src/main/resources/com/github/yassine/soxychains/subsystem/service/dns/config/pdnsd/pdnsd.supervisor.conf
[program:pdnsd]
command = pdnsd --daemon --config-file=/data/pdnsd/config/pdnsd.conf
autorestart = true