# src/main/resources/com/github/yassine/soxychains/subsystem/layer/spi/ovpn/config/dante/dante.supervisor.conf
[program:dante]
command = sockd -f /etc/config/dante/dante.conf
autorestart  = true
startretries = 10000