# src/main/resources/com/github/yassine/soxychains/subsystem/service/consul/config/consul/consul.supervisor.conf
[program:consul]
command = consul agent -enable-script-checks -dev -config-dir /etc/config/consul
autorestart = true
