# src/main/resources/com/github/yassine/soxychains/subsystem/service/gobetween/config/gobetween/gobetween.supervisor.conf
[program:gobetween]
command = gobetween -c /etc/config/gobetween/gobetween.toml
autorestart = true