# src/main/resources/com/github/yassine/soxychains/subsystem/layer/spi/ovpn/config/openvpn/openvpn.supervisor.conf
[program:openvpn]
command = /etc/config/openvpn/openvpn.sh